# schedule/com.snowflake.dailyupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.snowflake.dailyupdate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ayushi/Documents/snowflake-daily-update/venv/bin/python</string>
        <string>/Users/ayushi/Documents/snowflake-daily-update/daily_financial_update.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ayushi/Documents/snowflake-daily-update</string>
    <key>StandardOutPath</key>
    <string>/Users/ayushi/Documents/snowflake-daily-update/logs/daily_update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ayushi/Documents/snowflake-daily-update/logs/daily_update.log</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
